# supervisor.conf
[supervisord]

[program:nessus]
command			= /usr/bin/nessus_start.sh
autorestart		= true
priority		= 1